# deploy/dashboard.service
[Unit]
Description=eiDOS dashboard (read-only status server)
After=eidos.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/eidos
ExecStart=/usr/bin/python3 dashboard.py --config config.toml
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
