# deploy/llama-server.service
[Unit]
Description=llama.cpp inference server
After=network.target

[Service]
Type=simple
User=pi
ExecStart=/usr/local/bin/llama-server \
    -m /home/pi/models/qwen3-1.7b-q4_k_m.gguf \
    --host 127.0.0.1 \
    --port 8080 \
    -ngl 0 \
    -c 8192 \
    -np 1 \
    -t 4 \
    --reasoning-budget 0
Restart=always
RestartSec=10
Environment=LLAMA_LOG_LEVEL=warn

[Install]
WantedBy=multi-user.target
